# helpers/service/update_raid-daily.timer
[Unit]
Description=Update raid data daily sync at midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
